# packer/files/puma.service
[Unit]
Description=puma
After=multi-user.target

[Service]
Type=forking
ExecStart=/home/appuser/start-puma.sh

[Install]
WantedBy=multi-user.target

